# deploy/systemd/smith-delivery.target
[Unit]
# Groups the Smith basic-delivery worker pool + webhook trigger.
#
# Repo-tracked TEMPLATE. The install script (deploy/install.sh) copies this into
# ~/.config/systemd/user/ verbatim. Each unit declares
# `WantedBy=smith-delivery.target`, so enabling them wires them under this target;
# `systemctl --user start smith-delivery.target` then brings the whole pool up,
# and `stop` takes it down. Forgejo + the host-mode runner stay independent in
# the existing forgejo.service.
Description=Smith basic-delivery worker pool (architect + engineer + mechanical + trigger)
Documentation=https://git.ekanayaka.io/ai/smith/src/branch/main/docs/how-to/run-local-delivery.md
After=forgejo.service
Wants=forgejo.service

[Install]
WantedBy=default.target
