# deploy/systemd/smith-architect.service
[Unit]
# Smith basic-delivery architect role worker.
#
# Repo-tracked TEMPLATE. The install script (deploy/install.sh) copies this into
# ~/.config/systemd/user/ verbatim — there are no @PLACEHOLDER@ substitutions.
# Everything machine-specific is read at runtime:
#   - operator knobs (BASE_URL, OWNER, NAME, POLL_MS, …) from smith.env,
#   - per-role secrets (TEMPER_FORGEJO_*_ARCHITECT) from secrets/roles.env,
#   - %h for your home.
# ExecStart is the smith-role-worker shim: systemd's `Environment=` cannot remap
# the suffixed roles.env vars onto the generic names temper-worker reads (it does
# not expand ${VAR}), so the shim does that at runtime — keeping secrets off argv.
#
# The architect claims one untriaged intake issue and rewrites its thin body into
# a complete, implementable code spec (the triage_intake_to_code set_body effect),
# using a READ-ONLY checkout at ~/.local/state/smith/architect/smith.
Description=Smith basic-delivery architect role worker
Documentation=https://git.ekanayaka.io/ai/smith/src/branch/main/docs/how-to/run-local-delivery.md
PartOf=smith-delivery.target
# Forgejo + the host-mode runner come up via the separate forgejo.service.
After=forgejo.service network-online.target
Wants=network-online.target

[Service]
Type=simple
# No secrets on argv: smith.env carries the non-secret knobs, roles.env the
# per-role tokens/passwords. Missing files fail the unit loudly.
EnvironmentFile=%h/.config/smith/smith.env
EnvironmentFile=%h/.config/smith/secrets/roles.env
# The architect's stable workspace is a pre-cloned, READ-ONLY worktree Temper
# reuses across ticks (fetch + checkout -B; never clone/wipe).
Environment=TEMPER_CODING_WORKSPACE_ROOT=%h/.local/state/smith/architect/smith
Environment=TEMPER_CODING_WORKSPACE_PUSH=0
Environment=TEMPER_CODING_WORKSPACE_PR_LABELS=implementation
Environment=TEMPER_WORKFLOW_FILE=%h/.config/smith/workflow.json
# Temper runs this via `sh -c`, which expands $HOME / $SMITH_CODING_AGENT_ARGS
# from the inherited environment at run time (systemd does not expand it here).
Environment="TEMPER_CODING_WORKSPACE_COMMAND=$HOME/.local/bin/smith-coding-agent $SMITH_CODING_AGENT_ARGS"
ExecStart=%h/.local/bin/smith-role-worker architect
Restart=on-failure
RestartSec=5

[Install]
WantedBy=smith-delivery.target
