# deploy/systemd/smith-trigger.service
[Unit]
# Smith basic-delivery webhook accelerator (local trigger).
#
# Repo-tracked TEMPLATE. The install script (deploy/install.sh) copies this into
# ~/.config/systemd/user/ verbatim.
#
# temper-trigger-forgejo receives Forgejo webhooks on a loopback port and fans
# each accepted delivery out to EVERY worker's wake socket (no role routing), so
# the pool reacts immediately instead of waiting for the next long poll. Forge +
# trigger are both local, so there is NO ssh tunnel. Polling remains the liveness
# backstop if a webhook is ever lost.
Description=Smith basic-delivery webhook trigger (wake accelerator)
Documentation=https://git.ekanayaka.io/ai/smith/src/branch/main/docs/how-to/run-local-delivery.md
PartOf=smith-delivery.target
After=forgejo.service network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/smith/smith.env
ExecStart=%h/.local/bin/temper-trigger-forgejo \
    --bind ${TRIGGER_BIND} \
    --webhook-secret-file %h/.config/smith/secrets/webhook-secret \
    --wake-secret-file %h/.config/smith/secrets/wake-secret \
    --wake-dir %h/.local/state/smith/wake
Restart=on-failure
RestartSec=5

[Install]
WantedBy=smith-delivery.target
